# etc/service/nginx/run
#!/bin/sh
echo "Starting NGINX"
/usr/sbin/nginx
